# jarvis/scripts/com.francescotomatis.jarvis-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.francescotomatis.jarvis-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ngrok</string>
        <string>http</string>
        <string>--url=__NGROK_URL__</string>
        <string>8765</string>
        <string>--log=stdout</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>15</integer>
    <key>StandardOutPath</key>
    <string>__ROOT__/data/jarvis-tunnel.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/data/jarvis-tunnel.log</string>
</dict>
</plist>
