# jarvis/scripts/com.francescotomatis.jarvis-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.francescotomatis.jarvis-watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>__ROOT__/.venv/bin/python</string>
        <string>-m</string>
        <string>src.ops.watchdog</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__ROOT__</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>__ROOT__/data/watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>__ROOT__/data/watchdog.log</string>
</dict>
</plist>
